# systemd/digiurban.service
[Unit]
Description=DigiUrban Application
Documentation=https://github.com/fernandinhomartins40/digiurban2
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=digiurban
Group=digiurban
WorkingDirectory=/opt/digiurban/current
ExecStart=/usr/bin/node apps/backend/dist/index.js
Restart=always
RestartSec=10
StartLimitInterval=60s
StartLimitBurst=3

# Environment variables
Environment=NODE_ENV=production
Environment=PORT=3003
EnvironmentFile=-/opt/digiurban/current/.env

# Logging
StandardOutput=append:/opt/digiurban/logs/digiurban.log
StandardError=append:/opt/digiurban/logs/digiurban-error.log
SyslogIdentifier=digiurban

# Health monitoring
TimeoutStartSec=60s
TimeoutStopSec=30s
KillMode=process
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target